# docker/rootfs/etc/s6-overlay/s6-rc.d/bt-audio-manager/finish
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Cleanup on service stop
# ==============================================================================

bashio::log.info "Bluetooth Audio Manager stopped."
