# bluetooth_audio_manager/rootfs/etc/s6-overlay/s6-rc.d/bt-audio-manager/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Start the Bluetooth Audio Manager service
# ==============================================================================

bashio::log.info "Starting Bluetooth Audio Manager v${BUILD_VERSION:-unknown} (${BUILD_SHA:-no-sha})..."
bashio::log.info "PYTHONPATH=${PYTHONPATH}"

# Diagnostics: verify Python can see the module
bashio::log.info "ls /usr/src/: $(ls /usr/src/ 2>&1)"
bashio::log.info "ls /usr/src/bt_audio_manager/: $(ls /usr/src/bt_audio_manager/ 2>&1)"
bashio::log.info "Python sys.path: $(python3 -c 'import sys; print(sys.path)' 2>&1)"
bashio::log.info "Python import test: $(python3 -c 'import bt_audio_manager; print("OK")' 2>&1)"

exec python3 -m bt_audio_manager
