# bluetooth_audio_manager/rootfs/etc/s6-overlay/s6-rc.d/bt-audio-manager/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Start the Bluetooth Audio Manager service
# ==============================================================================

bashio::log.info "Starting Bluetooth Audio Manager..."

exec python3 -m bt_audio_manager
